# python/crontab
10 21 * * * /usr/local/bin/python /app/update_sp500.py >> /app/cron_test.log 2>&1
